# crontab
0 * * * * Gt\Website\Github\Cache::refresh()